# soulseek.conf
[program:soulseek]
environment=HOME="/home/soulseek",DISPLAY=":1",USER="soulseek",TERM="xterm"
command=/usr/bin/soulseek
user=soulseek
autorestart=true
priority=400
